# config/uniwherebot.conf
[program:uniwherebot]
command = gunicorn app:app -b localhost:8000
directory = /home/www/bot_app
user = michele